# resume-generation-backend.service
[Unit]
Description=NodeJS ResumeGeneration Backend
After=network.target
 
[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/nodeapps/resume-generation-backend/
ExecStart=/usr/bin/npm run start
Restart=on-failure
 
[Install]
WantedBy=multi-user.target

